# frr/files/frr.init
#!/bin/sh /etc/rc.common

START=60
start() {
	/usr/sbin/frr.init start
}

stop() {
	/usr/sbin/frr.init stop
}
