# overthebox/files/etc/uci-defaults/1900-firewall
#!/bin/sh
# shellcheck disable=SC1091
# vim: set noexpandtab tabstop=4 shiftwidth=4 softtabstop=4 :

. /lib/functions.sh

_setup() {
	name=
	config_get name "$1" name

	case "$name" in
	lan|wan)
		uci -q rename "firewall.$1=$name"
		;;
	mptcp)
		interface_list=
		config_get interface_list "$1" network
		for interface in $interface_list; do
			uci -q del_list firewall.wan.network="$interface"
			uci -q add_list firewall.wan.network="$interface"
		done
		uci -q delete firewall."$1"
		;;
	esac
}

_delete_mptcp_forwarding() {
	src=
	dest=
	config_get src "$1" src
	config_get dest "$1" dest

	if [ "$src" = "mptcp" ] || [ "$dest" = "mptcp" ]; then
		uci -q delete "firewall.$1"
	fi
}

_update_network_list() {
	cfg="$1"
	interface_list=
	config_get interface_list "$cfg" network

	# Delete this interface
	uci -q delete "firewall.$cfg.network"

	# Add each interfaces individually
	for interface in $interface_list; do
		uci -q del_list "firewall.$cfg.network=$interface"
		uci -q add_list "firewall.$cfg.network=$interface"
	done
}

_lowercase_proto() {
	cfg="$1"
	proto=
	config_get proto "$cfg" proto

	# shellcheck disable=SC2018,2019
	lcProto=$(echo "$proto" | tr 'A-Z' 'a-z')

	[ "$proto" != "$lcProto" ] && uci -q set "firewall.$cfg.proto=$lcProto"
}

config_load firewall
config_foreach _update_network_list zone
config_foreach _setup zone
config_foreach _delete_mptcp_forwarding forwarding
config_foreach _lowercase_proto rule
config_foreach _lowercase_proto redirect

uci -q batch <<-EOF
set firewall.lan.mtu_fix=1
del_list firewall.lan.network=if0
del_list firewall.wan.network=if0
add_list firewall.wan.network=if0
set firewall.allow_dhcp_request=rule
set firewall.allow_dhcp_request.name=Allow-DHCP-Request
set firewall.allow_dhcp_request.src=wan
set firewall.allow_dhcp_request.proto=udp
set firewall.allow_dhcp_request.dest_port=67
set firewall.allow_dhcp_request.target=ACCEPT
set firewall.allow_dhcp_request.family=ipv4
EOF
