# overthebox/files/etc/uci-defaults/1920-otb-0.6-to-0.7
#!/bin/sh
# shellcheck disable=SC1091,SC2039
# vim: set noexpandtab tabstop=4 shiftwidth=4 softtabstop=4 :

mode="$(uci -q get glorytun.app.mode)"

case "$mode" in
	glorytun_hybrid)
		logger -t update -p 6 "Update glorytun from xtun0"
		key="$(uci -q get glorytun.xtun0.key)"
		uci -q batch <<-EOF
		set glorytun.tun0='mud'
		set glorytun.tun0.key='$key'
		delete glorytun.xtun0
		delete glorytun.app
		EOF
		;;
	glorytun_mud)
		logger -t update -p 6 "Delete tun0.mtu"
		uci -q batch <<-EOF
		delete glorytun.tun0.mtu
		delete glorytun.tun0.mtu_auto
		delete glorytun.app
		EOF
		;;
esac

ipaddr="$(uci -q get network.xtun0.ipaddr)"
gateway="$(uci -q get network.xtun0.gateway)"

if [ "$ipaddr" ] && [ "$gateway" ]; then
	logger -t update -p 6 "Remove xtun0"
	uci -q batch <<-EOF
	set network.tun0.ipaddr='$ipaddr'
	set network.tun0.gateway='$gateway'
	delete network.tun0.mtu
	delete network.xtun0
	EOF
fi

if [ "$(uci -q get "network.if0")" = interface ]; then
	logger -t update -p 6 "Rename if0 to wan"
	uci -q batch <<-EOF
	rename network.if0=wan
	rename network.if0_dev=wan_dev
	set network.wan.ifname=wan
	set network.wan_dev.name=wan
	del_list firewall.wan.network=if0
	del_list firewall.wan.network=wan
	add_list firewall.wan.network=wan
	EOF
fi

# delete old stuff
sed -i '/169.254.254.1 api/d; /213.186.33.11 provisionning.overthebox.net/d' /etc/hosts
uci -q delete "dhcp.api_hostrecord"
rm -f /etc/sysctl.d/shadowsocks.conf
rm -f /etc/config/shadowsocks
rm -f /etc/config/otb-tracker
