# overthebox/files/etc/uci-defaults/otb-ucitrack.defaults
#!/bin/sh
# vim: set noexpandtab tabstop=4 shiftwidth=4 softtabstop=4 :

# TODO: do we still need this ?
touch /etc/config/ucitrack
