# overthebox/files/etc/uci-defaults/otb-firewall.defaults
#!/bin/sh
# vim: set noexpandtab tabstop=4 shiftwidth=4 softtabstop=4 :

# fix old conf with unnamed mptcp zone

name=$(uci -qX show firewall | grep "name='mptcp'" | cut -d. -f2)

[ "$name" = mptcp ] && exit 0

uci -q batch <<EOF
rename "firewall.$name='mptcp'"
commit firewall
EOF

exit 0
