# overthebox/files/etc/uci-defaults/otb-crontab.defaults
#!/bin/sh
# vim: set noexpandtab tabstop=4 shiftwidth=4 softtabstop=4 :

FILE="/etc/crontabs/root"

[ -f "$FILE" ] || exit 0

sed -i '/overthebox_download_recovery/d' "$FILE"
sed -i '/overthebox_checkup/d' "$FILE"

/etc/init.d/cron restart

exit 0
