# overthebox/files/etc/uci-defaults/otb-tracker.defaults
#!/bin/sh
# vim: set noexpandtab tabstop=4 shiftwidth=4 softtabstop=4 :

uci -q batch <<EOF
delete otb-tracker.defaults
set otb-tracker.defaults=defaults
add_list otb-tracker.defaults.hosts='213.186.33.99'
add_list otb-tracker.defaults.hosts='1.1.1.1'
add_list otb-tracker.defaults.hosts='8.8.8.8'
set otb-tracker.defaults.interval=1
EOF
