# overthebox/files/etc/uci-defaults/otb-uhttpd.defaults
#!/bin/sh
# vim: set noexpandtab tabstop=4 shiftwidth=4 softtabstop=4 :

# Do not redirect http to https, allow both
# overthebox.ovh will redirect to http for easier configuration
uci set uhttpd.main.redirect_https='0'
uci commit uhttpd
