# overthebox/files/etc/uci-defaults/otb-key.defaults
#!/bin/sh
# vim: set noexpandtab tabstop=4 shiftwidth=4 softtabstop=4 :

OTB_KEY_FILE=/etc/opkg/keys/6bc78405bb9b7cc3

[ -d "$OTB_KEY_FILE" ] && exit

cat > "$OTB_KEY_FILE" <<EOF
untrusted comment: OTB build key
RWRrx4QFu5t8wxTWDz+ErjM4vJhhhPzPvEJqaheYtgQXZuFSa+ndn91U
EOF
