# overthebox/files/etc/uci-defaults/otb-power.defaults
#!/bin/sh
# vim: set noexpandtab tabstop=4 shiftwidth=4 softtabstop=4 :

RUNTIME_POWER_SCRIPT=/etc/rc.button/power
CUSTOM_POWER_SCRIPT=/tmp/overthebox-power-button.sh

if [ -f ${CUSTOM_POWER_SCRIPT} ]; then
	mv ${CUSTOM_POWER_SCRIPT} ${RUNTIME_POWER_SCRIPT}
fi

exit 0
