# uqmi/files/usr/bin/uqmid
#!/bin/sh /etc/rc.common
# shellcheck disable=SC2034
# vim: set noexpandtab tabstop=4 shiftwidth=4 softtabstop=4 :
# Copyright (C) 2024 ovhcloud.com

_log() {
    logger -p daemon.info -t "LTE-Manager" "$@"
}

_err() {
    logger -p daemon.err -t "LTE-Manager" "$@"
}

# Function to start LTE interface
start_lte_interface() {
    # Get network interface name
    interface=$(qmicli --device=/dev/cdc-wdm0 --device-open-proxy --get-wwan-iface)
    if ifconfig | grep -q "${interface}"; then
        _log "LTE interface is already up."
        return 0
    fi

    _log "Bringing up LTE interface"
    ifup lte || {
        _err "Failed to bring up LTE interface"
        return 1
    }
}

# LTE device
device=$(ls /dev/* | grep -o '/dev/cdc-wdm[0-9]*' | uniq)
[ "${device}" ] || {
    _err "failed to retrieve LTE device"
    exit 1
}

# Main loop
while true; do
    # Check LTE connection status
    status=$(timeout 10 uqmi -d "$device" --get-data-status|tr -d '"')
    if [ "$?" -ne 0 ] || [ "${status}" != "connected" ]; then
        _log "LTE connection is down. Attempting to reconnect..."
        start_lte_interface || continue
        # Attempt to reconnect using uqmi
        start=$(timeout 10 uqmi -d "$device" --start-network| tr -d '"')
        if [ "$?" -ne 0 ]; then
            _err "Failed to reconnect using uqmi"
            # If reconnection fails, retry after a short delay
            sleep 5
            continue
        fi
    fi

    _log "LTE connection is up."
    # Sleep for a while before checking again
    sleep 30
done
