# mwan3otb/files/etc/hotplug.d/net/04-mwan3
#!/bin/sh
# vim: set noexpandtab tabstop=4 shiftwidth=4 softtabstop=4 :
/usr/bin/lua -e "require('overthebox').update_confmwan()"
