# mwan3otb/files/etc/hotplug.d/net/04-mwan3
#!/bin/sh
/usr/bin/lua -e "require('overthebox').update_confmwan()"
